# system_files/usr/lib/systemd/system/touch-etc-niri-local.kdl.service

[Unit]
Description = Ensure /etc/niri/local.kdl exists

[Service]
Type = oneshot
ExecStart = touch /etc/niri/local.kdl

[Install]
WantedBy = multi-user.target
